# services/telegrambot.service
[Unit]
Description=Telegram Bot
After=network.target

[Service]
User=solros
Group=solros
Type=simple
WorkingDirectory=/home/solros/defichain_price_alarm
ExecStart=python /home/solros/defichain_price_alarm/price_alarm_bot.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target